# dot_config/systemd/user/protonmail-bridge.service
[Unit]
Description=ProtonMail Bridge (headless IMAP/SMTP proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/protonmail-bridge --noninteractive
Restart=on-failure
RestartSec=10s
Environment=PASSWORD_STORE_DIR=%h/.local/share/password-store

[Install]
WantedBy=graphical-session.target
