# dot_config/systemd/user/wob.service
[Unit]
Description=Wayland overlay bar (volume / brightness OSD)
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStartPre=/usr/bin/sh -c 'rm -f "$XDG_RUNTIME_DIR/wob.sock" && mkfifo "$XDG_RUNTIME_DIR/wob.sock"'
ExecStart=/usr/bin/sh -c 'exec tail -f "$XDG_RUNTIME_DIR/wob.sock" | %h/.nix-profile/bin/wob'
ExecStopPost=/usr/bin/rm -f %t/wob.sock
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=sway-session.target
