# dot_config/systemd/user/zellij-inhibit-suspend.path
[Unit]
Description=Activate suspend inhibitor whenever zellij has a live session

[Path]
# %t expands to $XDG_RUNTIME_DIR (typically /run/user/$UID); zellij keeps
# its per-version session sockets under this directory. Whenever the dir
# transitions from empty to non-empty, the service is (re)activated.
DirectoryNotEmpty=%t/zellij
Unit=zellij-inhibit-suspend.service

[Install]
WantedBy=default.target
