# home/.config/systemd/user/gpg-agent-ssh.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=%t/gnupg/d.199epr64wmzkrnk8u8qgricf/S.gpg-agent.ssh