# home/.config/systemd/user/gpg-agent-extra.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=%t/gnupg/d.hmaqciuk8y8ye3gwt9b6eth1/S.gpg-agent.extra