# home/.config/systemd/user/vdirsyncer.service.d/override.conf
[Service]
Environment="GNUPGHOME=%h/.local/share/gnupg"
Environment="PASSWORD_STORE_DIR=%h/.local/share/password-store"
