# etc/systemd/system/getty@.service.d/blank.conf
[Service]
ExecStartPre=/bin/sh -c 'setterm --blank 5 --powerdown 10 >/dev/tty%I'
